# indiwebmanager.service
[Unit]
Description=INDI Web Manager
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python /home/pi/servermanager/drivermanager.py

[Install]
WantedBy=multi-user.target
